# transcription-system/scripts/systemd/transcription-backend.service
[Unit]
Description=Transcription System Backend
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/transcription-system/transcription-system/backend
Environment="NODE_ENV=production"
Environment="PORT=5000"
ExecStart=/usr/bin/node dist/server.js
Restart=always
RestartSec=10
StandardOutput=append:/var/log/transcription-backend.log
StandardError=append:/var/log/transcription-backend.log

[Install]
WantedBy=multi-user.target